# deploy-system-unified/roles/storage/backup/restic/templates/restic-backup.timer.j2
[Unit]
Description=Daily Restic Backup Timer

[Timer]
OnCalendar={{ restic_schedule_on_calendar }}
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
