# projects/deploy-system-unified/roles/containers/monitoring/templates/grafana.container.j2
[Unit]
Description=Grafana Container
After=network-online.target

[Container]
Image={{ monitoring_grafana_image }}
Pod={{ monitoring_pod_name }}
EnvironmentFile={{ containers_secrets_dir }}/grafana.env
Volume={{ monitoring_root_dir }}/grafana:/var/lib/grafana:Z

[Install]
WantedBy=multi-user.target
