# roles/security/file_integrity/templates/aide-boot-check.service.j2
{# =============================================================================
# Audit Event Identifier: DSU-TPL-900021
# Template Type: file_integrity Configuration
# Last Updated: 2026-02-28
# Version: 1.0
# ============================================================================= #}
[Unit]
Description=AIDE File Integrity Check on Boot
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/aide --check
StandardOutput=journal+console
StandardError=journal+console
User=root
Group=root

[Install]
WantedBy=multi-user.target